# deploy/menuvi.service
[Unit]
Description=MenuVi Gunicorn Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/menuvi
EnvironmentFile=/var/www/menuvi/.env
ExecStart=/var/www/menuvi/.venv/bin/gunicorn "menuvi:create_app()" \
    --bind 127.0.0.1:8010 \
    --workers 2 \
    --timeout 30 \
    --access-logfile /var/log/menuvi/access.log \
    --error-logfile /var/log/menuvi/error.log
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
